# debian/openfpc-node.postrm
#! /bin/sh
# postrm script for openfpc-node
#
# see: dh_installdeb(1)
# see /usr/doc/packaging-manual/

set -e

case "$1" in
       purge|disappear)
                # Remove logfiles
                rm -rf /var/tmp/openfpc/

                # Remove configuration
                if [ -e /etc/openfpc/openfpc-default.conf ] ; then
                        rm /etc/openfpc/openfpc-example-node.conf
                fi

                # Remove configuration dir
                [ -d /etc/openfpc ] && rmdir --ignore-fail-on-non-empty /etc/openfpc

                # Remove user/group
                if getent passwd | grep -q "^openfpc:"; then
                        userdel openfpc 2>/dev/null || true
                fi
                if getent group | grep -q "^openfpc:" ; then
                        delgroup --only-if-empty openfpc 2>/dev/null || true
                fi
                ;;
       remove) 
                update-rc.d -f openfpc-daemonlogger remove >/dev/null
                update-rc.d -f openfpc-queued remove >/dev/null 
                update-rc.d -f openfpc-cxtracker remove >/dev/null
                update-rc.d -f openfpc-cx2db remove >/dev/null
	;;
       upgrade|failed-upgrade|abort-install|abort-upgrade)
                # nothing
                # We may not delete the user openfpc, as there may be
                # files owned by it in /var/tmp/openfpc and /etc/openfpc.
                ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#
